# etc/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:sphere-go-state-service]
command=/app/sphere-go-state-service
redirect_stderr=true