# etc/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:sphere-go-state-service]
command=/go/src/github.com/ninjablocks/sphere-go-state-service/sphere-go-state-service
redirect_stderr=true